# backend/server_manager/server-manager.service
[Unit]
Description=Hiddify Server Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/hiddify-app-custom/backend
EnvironmentFile=-/root/hiddify-app-custom/backend/server_manager/.env
ExecStart=/root/hiddify-app-custom/backend/server_manager/.venv/bin/python -m server_manager
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
